# monitor/genesis-monitor.conf
[program:genesis-monitor]
command=python3 /opt/genesis-monitor/monitor/monitor_worker.py
directory=/opt/genesis-monitor/monitor
autostart=true
autorestart=true
startsecs=10
startretries=5
user=root
redirect_stderr=true
stdout_logfile=/var/log/genesis-monitor.out.log
stderr_logfile=/var/log/genesis-monitor.err.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PATH="/usr/local/bin:/usr/bin:/bin"
